# roles/hadoop-base/templates/etc/profile.d/kdp.sh.j2

for envfile in /etc/kdp/env/*
do
    source $envfile
done
