# roles/tinydoop/files/etc/profile.d/99-airflow.sh
export AIRFLOW_CONFIG=/etc/airflow/airflow.cfg
export AIRFLOW_HOME=/etc/airflow/
export PATH="$PATH:/opt/airflow/bin"
